# startup/vdr-pre.conf
# Erledigt Vorbereitungen für den VDR

start on startup

exec /usr/sbin/vdr-preupstart
