# startup/vdr-pre.conf
# Erledigt Vorbereitungen für den VDR

task

start on filesystem

exec /usr/sbin/vdr-preupstart
