# startup/yavdr-startup/udev/40-dvb-devices.rules
# restart vdr if dvb device added after vdr is running
KERNEL=="dvb?.frontend?", ACTION=="add", RUN+="/lib/udev/dvbmon_helper"
