# misc/yavdr-dvbmon/yavdr-dvbmon.conf
# yavdr-dvbmon - yaVDR DVB device monitor
#
# rsyslog is an enhanced multi-threaded replacement for the traditional
# syslog daemon, logging messages from applications

description	"yaVDR DVB device monitor"

start on filesystem
stop on runlevel [06]


exec yavdr-dvbmon 
