# misc/dvbmon/etc.init.dvb-first.conf
start on local-filesystems

task

exec /usr/bin/dvbmon
