# misc/dvbmon/usr.udev.rules.d.40-dvb-devices.rules
# check if we are finished with devices
KERNEL=="dvb?.frontend?", ACTION=="add", RUN+="/usr/bin/dvbmon"

