# vdr-addon-avahi-mounter/avahi-mounter.conf

author "Arno Esser <arno.esser@gmx.de>"
description "Use avahi to mount resources from network"

start on filesystem and net-device-up IFACE!=lo
stop on runlevel [016]

 respawn 

script

if [ -e /etc/default/locale ]; then 
. /etc/default/locale
fi 

if [ -n "$LANG" ]; then 
   export LANG
fi 

exec /usr/bin/avahi-mounter
end script

