# startup/templates/etc/init/x.conf
start on (filesystem
          and started dbus
          and (graphics-device-added fb0 PRIMARY_DEVICE_FOR_DISPLAY=1
               or drm-device-added card0 PRIMARY_DEVICE_FOR_DISPLAY=1
               or stopped udevtrigger))
stop on runlevel [016]

respawn

env XORGCONFIG=/etc/X11/xorg.conf.yavdr
export XORGCONFIG

exec /usr/bin/X :1 vt9

post-start script
while ! DISPLAY=:1 xli -onroot -center -border black /usr/share/yavdr/images/yavdr_logo.png ; do
sleep 1
done
end script
