# base/yavdr-base/etc/init/dvb-first.conf
start on (startup and net-device-up IFACE!=lo)

task

exec /lib/udev/dvbmon_helper
