# base/yavdr-base/etc/init/dvb-last.conf
start on stopped udev-finish

task

script
# wait politely to let really all hardware finish
sleep 30

# check for special handling. e.g. sundtek
if [ -d /etc/yavdr/init/dvb-last.d ]; then
    run-parts -v /etc/yavdr/init/dvb-last.d
fi

# now fire the vdr start as last fall back 
/sbin/initctl emit --no-wait dvb-devices-complete

# and save everything and log out result
DEVCNT="`ls /dev/dvb/adapter?/frontend? -1 | wc -l`" || /bin/true
logger -t "dvbmon" "Saving number of devices ($DEVCNT) after boot" || /bin/true
dbset system.vdr.dvbmon.device_cnt="$DEVCNT" || /bin/true

if [ "$DVBMON" = "1" ]; then
   logger -t "dvbmon" "Saving new number of devices ($DEVCNT)" || /bin/true
fi

env > /tmp/dvb-last
ls -la /tmp/ >> /tmp/dvb-last

if [ -e /tmp/.dvb-devices-complete ]; then
   logger -t "dvbmon" "vdr started by udev allready" || /bin/true
   rm -f /tmp/.dvb-devices-complete
fi

if [ -e /tmp/.dvb-first ]; then
   logger -t "dvbmon" "vdr started by dvb-first directly after mount" || /bin/true
   rm -f /tmp/.dvb-first
fi
end script

