# yavdr-hardware/mdm166a/yavdr-hardware-mdm166a-debian/debian/postinst
#!/bin/sh

case "$1" in
    configure)       
       /etc/init.d/LCDd stop || /bin/true
       /usr/sbin/update-rc.d -f LCDd remove || /bin/true
       [ -f /var/lib/yavdrdb.hdf ] && signal-event post-upgrade
       ;;

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

exit 0

#DEBHELPER#
